# worker/crontab
0 2 * * * root node /app/sendEmail.js >> /var/log/cron.log 2>&1

